# sv/elasticsearch/run
#!/bin/sh
exec 2>&1
exec /opt/elasticsearch-0.90.5/bin/elasticsearch -f
