# sv/elasticsearch/run
#!/bin/sh
exec 2>&1
exec /opt/elasticsearch/bin/elasticsearch -f